# remaster/mcs/mcs-oemsystem-default
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

. /sbin/mcs-functions.sh

MCS_CONFIGURED="/etc/.mcs-oemsystem-default-configured"

depend() {
	need net
	use dns
	before 389-ds 389-ds-snmp 389-admin
	before mysql jboss-bin-4.2 postfix ejabberd
}

start() {
	if [ -f "${MCS_CONFIGURED}" ]; then
		return 0
	fi
	einfo "Setting up MailWare Collaboration Suite"
	setup_fds
	rc="$?"
	eend "${rc}"
	if [ "${rc}" = "0" ]; then
		touch "${MCS_CONFIGURED}"
	fi
	return "${rc}"
}
